# cheatsheets/epmd.ini
[program:epmd]
command=epmd -port {{ epmd_port }}
user={{ epmd_user }}
autorestart=true
stdout_logfile={{ stack_log_path }}/epmd/epmd.stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile={{ stack_log_path }}/epmd/erlang.stderr.log
stderr_logfile_maxbytes=10MB
priority=300
