# cheatsheets/gunicorn.ini
[program:project_{{ django_env }}_gunicorn]
command={{ django_path }}/startup.sh {{ django_env }} {% if gunicorn_reload %}--reload{% endif %}
directory={{ django_path }}
user={{ user }}
autorestart=true
killasgroup=true
stopasgroup=true
stdout_logfile={{ stack_log_path }}/gunicorn/gunicorn.stdout.log
stdout_logfile_maxbytes=80MB
stderr_logfile={{ stack_log_path }}/gunicorn/gunicorn.stderr.log
stderr_logfile_maxbytes=80MB
priority=1
