# cheatsheets/rabbitmq.ini
[program:rabbitmq]
command={{ rabbitmq_sbin_path }}/rabbitmq-server
directory={{ rabbitmq_sbin_path }}
user={{ user }}
autorestart=true
stdout_logfile={{ stack_log_path }}/rabbitmq/rabbitmq.stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile={{ stack_log_path }}/rabbit/rabbit.stderr.log
stderr_logfile_maxbytes=10MB
priority=400
